# res/org.filezilla-project.filezilla-server.service.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.filezilla-project.filezilla-server.service</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/@PACKAGE_NAME@.app/Contents/MacOS/filezilla-server</string>
		<string>--config-dir=/Library/Preferences/org.filezilla-project.filezilla-server.service</string>
	</array>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>wheel</string>
	<key>Disabled</key>
	<false/>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
